# packaging/apt/bzero-agent.conf
description     "BastionZero Agent"
author          "BastionZero Inc."

start on runlevel [2345]
stop on runlevel [!2345]

respawn 

exec /usr/bin/$BZERO_PKG_NAME
